# system_files/usr/lib/systemd/system/ublueos-setup-docker.service
[Unit]
Description=Setup docker group
Before=multi-user.target

[Service]
ExecStart=/usr/libexec/ublueos-setup-docker

[Install]
WantedBy=basic.target